# .config/nvim/.skeleton_systemd.service
[Unit]
Description=

[Service]
Type=
WorkingDirectory=
ExecStartPre=
ExecStart=/bin/bash -c ''
ExecStop=

[Install]
WantedBy=

# man systemd.service
